# systemd/pop-updater.service
[Unit]
Description=Daily system and app update (pop-updater)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pop-updater
StandardOutput=journal
StandardError=journal
EnvironmentFile=/etc/pop-updater/secrets.conf
