# systemd/pop-updater.timer
[Unit]
Description=Run pop-updater daily at 07:30 WIB

[Timer]
OnCalendar=*-*-* 07:30:00
Persistent=true

[Install]
WantedBy=timers.target
